# deploy/workflow-autoheal.service
# Workflow autoheal — restart unhealthy workflow containers.
#
# docs/specs/daemon-liveness-watchdog.md §1 acceptance: a failing
# healthcheck must lead to a restart. Docker Compose does not restart
# unhealthy containers on its own (restart policies act on exit, not
# health), so this host-side oneshot closes the loop without granting
# any container docker.sock access.
#
# Triggered by workflow-autoheal.timer (every 2 minutes).

[Unit]
Description=Restart unhealthy workflow containers (liveness watchdog)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'unhealthy=$(docker ps -q --filter "health=unhealthy" --filter "label=org.workflow.component"); if [ -n "$unhealthy" ]; then for c in $unhealthy; do name=$(docker inspect -f "{{.Name}}" "$c"); echo "autoheal: restarting unhealthy $name"; docker restart "$c"; done; else exit 0; fi'
SyslogIdentifier=workflow-autoheal
