# deploy/workflow-disk-watch.service
# Workflow disk pressure watchdog — daily /var/lib/docker space check.
#
# Fires daily. If /var/lib/docker exceeds DISK_WARN_PCT (default 80%),
# opens a `disk-pressure` GitHub Issue via GITHUB_TOKEN so the host
# gets an email alert before Docker builds start failing.
#
# Mirrors the GH-issue alarm pattern from scripts/watchdog.py.
#
# Environment variables sourced from /etc/workflow/env:
#   GITHUB_TOKEN       — fine-grained PAT or GITHUB_TOKEN with issues:write
#   GITHUB_REPOSITORY  — owner/repo  (default: Jonnyton/Workflow)
#   DISK_WARN_PCT      — integer threshold, default 80
#   DISK_WATCH_PATH    — path to check, default /var/lib/docker
#
# Install:
#   sudo cp deploy/workflow-disk-watch.service /etc/systemd/system/
#   sudo cp deploy/workflow-disk-watch.timer    /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now workflow-disk-watch.timer

[Unit]
Description=Workflow disk pressure watchdog (daily /var/lib/docker alert)
Documentation=https://github.com/jfarnsworth/workflow/blob/main/docs/ops/log-aggregation-runbook.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/workflow/env
ExecStart=/usr/bin/python3 /opt/workflow/scripts/disk_watch.py
TimeoutStartSec=60s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=workflow-disk-watch
