# deploy/workflow-watchdog.service
# Workflow daemon watchdog — one-shot probe + restart.
#
# Self-host migration Row L per
# docs/exec-plans/active/2026-04-20-selfhost-uptime-migration.md.
#
# Triggered by workflow-watchdog.timer every 2 min. Runs a single
# probe-tick via scripts/watchdog.py and exits. State persists at
# /var/lib/workflow-watchdog/state.json across ticks.
#
# Install (handled by hetzner-bootstrap.sh):
#   sudo cp deploy/workflow-watchdog.service /etc/systemd/system/
#   sudo cp deploy/workflow-watchdog.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now workflow-watchdog.timer

[Unit]
Description=Workflow daemon watchdog (probe + restart on 3 consecutive reds)
Documentation=https://github.com/Jonnyton/Workflow/blob/main/deploy/HETZNER-DEPLOY.md
After=workflow-daemon.service
# Don't race the daemon's own startup at boot. Give daemon 60s to bind.

[Service]
Type=oneshot
User=workflow
Group=workflow
WorkingDirectory=/opt/workflow

# /var/lib/workflow-watchdog needs to be writable by the workflow user.
# systemd creates + chowns it on first start.
StateDirectory=workflow-watchdog
StateDirectoryMode=0750

# Need sudo to issue systemctl restart. The workflow user gets a
# tightly-scoped sudoers rule installed by the bootstrap (see
# HETZNER-DEPLOY.md): workflow ALL=(root) NOPASSWD:/usr/bin/systemctl
# restart workflow-daemon.service
#
# scripts/watchdog.py calls `systemctl restart workflow-daemon.service`
# when threshold is crossed. That command hits the sudoers allowlist.

ExecStart=/usr/bin/python3 /opt/workflow/scripts/watchdog.py

# One-shot exit; timer handles cadence. Fail fast if script errors.
TimeoutStartSec=60s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=workflow-watchdog
