# deploy/workflow-disk-watch.timer
# Workflow disk pressure watchdog timer — fires hourly.
#
# Self-host migration Row K (ops hardening) per
# docs/exec-plans/active/2026-04-20-selfhost-uptime-migration.md.
#
# Cadence bumped daily→hourly 2026-04-24 (BUG-023 Phase 4): daily was
# too slow to catch the silent disk-fill outage class. Hourly lets
# auto-prune intercept before compose crash-loops.

[Unit]
Description=Timer for workflow disk pressure watchdog (hourly)
Requires=workflow-disk-watch.service

[Timer]
# Hourly on :27 offset — avoids exact top-of-hour collision with other
# timers. OnBootSec runs shortly after boot so a fresh deploy sees
# pressure signals without waiting up to an hour.
OnBootSec=5min
OnUnitActiveSec=1h
OnCalendar=*:27
AccuracySec=5min
# Persistent=true so a missed check fires on next boot.
Persistent=true
Unit=workflow-disk-watch.service

[Install]
WantedBy=timers.target
