# deploy/workflow-ship-logs.timer
# Workflow log archiver timer — fires every hour.
#
# Self-host migration Row K per
# docs/exec-plans/active/2026-04-20-selfhost-uptime-migration.md.

[Unit]
Description=Timer for workflow log archiver (hourly)
Requires=workflow-ship-logs.service

[Timer]
# Every full hour (00:00, 01:00, ...). Log archiving is best-effort;
# a missed run recovers automatically via LOG_SINCE windowing.
OnCalendar=hourly
AccuracySec=5min
# Persistent=true so a missed ship (box rebooted mid-hour) fires on
# next boot to avoid gaps in the offsite archive.
Persistent=true
Unit=workflow-ship-logs.service

[Install]
WantedBy=timers.target
