# deploy/workflow-prune.timer
# Workflow docker image prune — weekly on Sunday 04:00 UTC.
#
# Runs after the nightly backup window (02:00–03:00) and before the
# tier-3 OSS clone (07:17). Persistent=true so a missed prune runs on
# next boot.

[Unit]
Description=Timer for workflow docker image prune
Requires=workflow-prune.service

[Timer]
OnCalendar=Sun 04:00 UTC
AccuracySec=5min
Persistent=true
Unit=workflow-prune.service

[Install]
WantedBy=timers.target
