# deploy/workflow-backup.service
# Workflow daemon state backup — oneshot snapshot + upload.
#
# Self-host migration Row J per
# docs/exec-plans/active/2026-04-20-selfhost-uptime-migration.md.
#
# Triggered by workflow-backup.timer nightly at 03:00 UTC. Invokes
# deploy/backup.sh to tar the workflow-data named volume + rclone-
# upload to Hetzner Storage Box. Retention pruned in the same run.
#
# Runs as root because:
#   1. /var/lib/docker/volumes/ is only root-readable by default.
#   2. `docker volume inspect` needs docker socket access — simpler
#      to grant via root than extend the workflow user's docker-group
#      membership across every ops primitive.
#
# Backup credentials live in /etc/workflow/env (the same env file the
# daemon uses). systemd sources it via EnvironmentFile; backup.sh
# reads STORAGEBOX_* from the process env.

[Unit]
Description=Workflow daemon state backup (nightly snapshot to Hetzner Storage Box)
Documentation=https://github.com/Jonnyton/Workflow/blob/main/deploy/HETZNER-DEPLOY.md
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/workflow/env
ExecStart=/bin/bash /opt/workflow/deploy/backup.sh
TimeoutStartSec=30min
# Don't kill mid-upload — backup.sh handles its own timeouts.

# PrivateTmp=yes so the rclone config written to $HOME is wiped
# on exit (defense-in-depth; backup.sh also shreds it).
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=workflow-backup
